# packages/lsmtpd_1.1.0-4_amd64/lib/systemd/system/lsmtpd.service
[Unit]
Description=LSMTPD - Lightweight SMTP Daemon
After=network.target
Documentation=https://lsmtp.NekoNik.org

[Service]
Type=simple
EnvironmentFile=/etc/default/lsmtpd
ExecStart=/usr/bin/lsmtpd
StandardOutput=append:/var/log/lsmtpd.log
StandardError=append:/var/log/lsmtpd.log
Restart=on-failure
User=root
Group=root
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
