# packages/lsmtpd_1.1.0-4_arm64/lib/systemd/system/lsmtpd.service
[Unit]
Description=LSMTPD - Lightweight SMTP Daemon
After=network.target
Documentation=https://lsmtp.NekoNik.org

[Service]
Type=simple
# Load environment from /etc/default/lsmtpd (created by postinst)
EnvironmentFile=/etc/default/lsmtpd
ExecStart=/usr/bin/lsmtpd
Restart=on-failure
# Run as dedicated unprivileged account
User=lsmtp
Group=lsmtp
LimitNOFILE=4096
# Do not place any sensitive credentials in the unit itself.

[Install]
WantedBy=multi-user.target
