# packages/lsmtpd_1.1.0-4_amd64/DEBIAN/postinst
#!/bin/sh
set -e

# Debconf helper
. /usr/share/debconf/confmodule || true

# If debconf is available, ask questions.
db_input high lsmtpd/description || true
db_input high lsmtpd/rust_log || true
db_input high lsmtpd/bind_address || true
db_input high lsmtpd/bind_port || true
db_input high lsmtpd/server_name || true
db_input high lsmtpd/amqp_host || true
db_input high lsmtpd/amqp_port || true
db_input high lsmtpd/amqp_username || true
db_input high lsmtpd/amqp_password || true
db_input high lsmtpd/amqp_vhost || true
db_input high lsmtpd/amqp_exchange || true
db_input high lsmtpd/amqp_routing_key || true
db_input high lsmtpd/amqp_buffer_size || true
db_input high lsmtpd/max_email_size || true

# Process the questions and populate debconf database
db_go || true

# Fetch values (fall back to defaults if db_get fails)
db_get lsmtpd/rust_log || true; RUST_LOG="${RET:-INFO}"
db_get lsmtpd/bind_address || true; BIND_ADDRESS="${RET:-0.0.0.0}"
db_get lsmtpd/bind_port || true; BIND_PORT="${RET:-25}"
db_get lsmtpd/server_name || true; SERVER_NAME="${RET:-nekonik.com}"
db_get lsmtpd/amqp_host || true; AMQP_HOST="${RET:-}"
db_get lsmtpd/amqp_port || true; AMQP_PORT="${RET:-}"
db_get lsmtpd/amqp_username || true; AMQP_USERNAME="${RET:-}"
db_get lsmtpd/amqp_password || true; AMQP_PASSWORD="${RET:-}"
db_get lsmtpd/amqp_vhost || true; AMQP_VHOST="${RET:-}"
db_get lsmtpd/amqp_exchange || true; AMQP_EXCHANGE="${RET:-}"
db_get lsmtpd/amqp_routing_key || true; AMQP_ROUTING_KEY="${RET:-}"
db_get lsmtpd/amqp_buffer_size || true; AMQP_BUFFER_SIZE="${RET:-10}"
db_get lsmtpd/max_email_size || true; MAX_EMAIL_SIZE_BYTES="${RET:-52428800}"


# Write env file (conservative permissions)
ENV_FILE="/etc/default/lsmtpd"
umask 077
mkdir -p /etc/default
cat > "$ENV_FILE" <<EOF
# /etc/default/lsmtpd - environment for lsmtpd (generated by package postinst)
RUST_LOG=lsmtpd=${RUST_LOG}
BIND_ADDRESS=${BIND_ADDRESS}
BIND_PORT=${BIND_PORT}
SERVER_NAME=${SERVER_NAME}
AMQP_HOST=${AMQP_HOST}
AMQP_PORT=${AMQP_PORT}
AMQP_USERNAME=${AMQP_USERNAME}
AMQP_PASSWORD='${AMQP_PASSWORD}'
AMQP_VHOST=${AMQP_VHOST}
AMQP_EXCHANGE=${AMQP_EXCHANGE}
AMQP_ROUTING_KEY=${AMQP_ROUTING_KEY}
AMQP_BUFFER_SIZE=${AMQP_BUFFER_SIZE}
MAX_EMAIL_SIZE_BYTES=${MAX_EMAIL_SIZE_BYTES}
EOF
chmod 0600 "$ENV_FILE"

# Ensure binary is executable and has reasonable permissions
if [ -f /usr/bin/lsmtpd ]; then
    chmod 0755 /usr/bin/lsmtpd || true
fi

# Reload systemd and enable service (do not start)
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
    # Enable so it's ready at boot, but DO NOT start now.
    systemctl enable lsmtpd.service >/dev/null 2>&1 || true
fi

# Inform debconf that configuration is done (suppress auto-start)
db_stop || true

exit 0
